# raon-vending.service
[Unit]
Description=RAON Vending Machine Kiosk Application
After=network.target display-manager.service

[Service]
Type=simple
User=pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
WorkingDirectory=/home/pi/raon-vending-rpi4
ExecStart=/usr/bin/python3 /home/pi/raon-vending-rpi4/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
