# deploy/raon-vending.service
[Unit]
Description=RAON Vending Machine Kiosk - Automatic Startup
After=network.target display-manager.service graphical.target

[Service]
Type=simple
User=pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
WorkingDirectory=/home/pi/raon-vending-rpi4
ExecStart=/bin/bash /home/pi/raon-vending-rpi4/deploy/start-kiosk-wayland.sh
Restart=on-failure
RestartSec=5
KillMode=process
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
