# raon-vending-rpi4-main/deploy/raon-vending.service
[Unit]
Description=RAON Vending Machine Kiosk - Automatic Startup
After=graphical.target

[Service]
Type=simple
# Run as pi user (default Raspberry Pi user)
User=pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
WorkingDirectory=/home/pi/raon-vending
ExecStart=/usr/bin/python3 main.py
# Restart on failure with 5 second delay
Restart=on-failure
RestartSec=5
# Kill the service on exit
KillMode=process

[Install]
# Start after graphical interface loads
WantedBy=graphical.target
